# infra/systemd/ecom-admin.service
[Unit]
Description=Ecom Admin (Next.js)
After=network.target
ConditionPathExists=/var/www/ecom/apps/admin
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/var/www/ecom/apps/admin
ExecStartPre=/usr/bin/test -x /usr/bin/node
ExecStartPre=/usr/bin/test -d /var/www/ecom/apps/admin
ExecStart=/usr/bin/node /var/www/ecom/apps/admin/render-start.js
Environment=NODE_ENV=production
Environment=NEXT_PUBLIC_API_BASE_URL=https://api.jeeey.com
Environment=NEXTAUTH_URL=https://admin.jeeey.com
Environment=NEXTAUTH_SECRET=${NEXTAUTH_SECRET}
Environment=PORT=3001
Restart=always
RestartSec=2
User=root
Group=root

[Install]
WantedBy=multi-user.target

