# infra/systemd/ecom-api.service
[Unit]
Description=Ecom API (Node)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/ecom/packages/api
ExecStart=/usr/bin/node -r dotenv/config dist/index.js
Environment=NODE_ENV=production
Environment=PORT=4000
Environment=API_FORCE_LISTEN=1
Environment=DOTENV_CONFIG_PATH=/var/www/ecom/.env.api
Restart=always
RestartSec=2
User=root
Group=root

[Install]
WantedBy=multi-user.target

