# infra/systemd/ecom-api.service
[Unit]
Description=Ecom API (Node)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/ecom
ExecStart=/usr/bin/node -r dotenv/config packages/api/dist/index.js
Environment=NODE_ENV=production
Environment=PORT=4000
Environment=API_FORCE_LISTEN=1
Environment=DOTENV_CONFIG_PATH=packages/api/.env
Restart=always
RestartSec=2
User=deploy
Group=deploy

[Install]
WantedBy=multi-user.target

