# infra/systemd/ecom-admin.service
[Unit]
Description=Ecom Admin (Next.js)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/ecom/apps/admin/.next/standalone
ExecStart=/usr/bin/node apps/admin/server.js
Environment=NODE_ENV=production
Environment=PORT=3001
Restart=always
RestartSec=2
User=root
Group=root

[Install]
WantedBy=multi-user.target

