# infra/systemd/ecom-api.service
[Unit]
Description=Ecom API (Node)
After=network.target
ConditionPathExists=/var/www/ecom/packages/api
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/var/www/ecom/packages/api
ExecStartPre=/usr/bin/test -x /usr/bin/node
ExecStartPre=/usr/bin/test -f /var/www/ecom/.env.api
ExecStart=/usr/bin/node -r dotenv/config dist/index.js
Environment=NODE_ENV=production
Environment=PORT=4000
Environment=API_FORCE_LISTEN=1
Environment=DOTENV_CONFIG_PATH=/var/www/ecom/.env.api
Environment=COOKIE_DOMAIN=.jeeey.com
Environment=NEXT_PUBLIC_ADMIN_URL=https://admin.jeeey.com
Environment=NEXT_PUBLIC_APP_URL=https://jeeey.com
Environment=CORS_ALLOW_ORIGINS=https://jeeey.com,https://www.jeeey.com,https://admin.jeeey.com,https://api.jeeey.com
Environment=API_RUN_ENSURE_SCHEMA=1
Environment=API_ALLOW_BOOTSTRAP=1
Environment=NLP_CONFIG_DIR=/var/www/ecom/config/nlp
Restart=always
RestartSec=2
User=root
Group=root

[Install]
WantedBy=multi-user.target

